# .config/systemd/user/theme.service
[Unit]
Description=Cambiar tema automáticamente según la hora

[Service]
Type=oneshot
ExecStart=%h/.scripts/services/theme.sh

